# redis/templates/monitor/monitor.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Run redis monitor perdiodically

[Timer]
OnBootSec=1m
OnUnitActiveSec=15s

[Install]
WantedBy=timers.target
